# Configuration/iems5722_app.conf
[program:gunicorn]
command = gunicorn -b localhost:8000 -w 4 A3_1155164941.wsgi
directory = /home/AndroidBackend/A3_1155164941
user = root
autostart = true
autorestart = true
stdout_logfile = /home/AndroidBackend/A3_1155164941/app.log
redirect_stderr = true

[program:celeryworker]
command=celery -A A3_1155164941 worker -l info
directory=/home/AndroidBackend/A3_1155164941
user=root
autostart=true
autorestart=true
stdout_logfile=/home/AndroidBackend/A3_1155164941/app.log
redirect_stderr=true